# quasi-senate/deploy/quasi-senate-solve.timer
[Unit]
Description=QUASI Senate Solve — every 4 hours (offset from draft)

[Timer]
OnCalendar=*-*-* 02/4:00:00
Persistent=true

[Install]
WantedBy=timers.target
